# crontab
SHELL=/bin/sh
0 3 * * * www cd /home/www/wereldontdooien && ./publish.sh
